# heartbeat.service
[Unit]
Description=Heartbeat para Healthchecks.io

[Service]
Type=oneshot
ExecStart=/usr/local/bin/heartbeat.sh

